# crates/kit/src/units/bcvk-copy-units.service
[Unit]
Description=Copy bcvk units for post-switch-root on systemd <256
DefaultDependencies=no
ConditionPathExists=/etc/initrd-release
Before=initrd-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
# On systemd <256, SMBIOS credentials with systemd.extra-unit.* are not processed
# natively. We work around this by embedding the units in the initramfs and copying
# them to /run/systemd/system/ before switch-root. The /run tmpfs is preserved
# across switch-root via MS_MOVE.
#
# Copy journal-stream unit and create dropin for sysinit.target, but only
# when the host has requested journal streaming (indicated by the presence of
# the SMBIOS credential that bcvk injects when --output=journal or
# --log-dir=journal=… is used).
ExecStart=/bin/sh -c 'test -e /run/credentials/@system/systemd.extra-unit.bcvk-journal-stream.service || exit 0; \
  mkdir -p /run/systemd/system/sysinit.target.wants; \
  cp /usr/lib/systemd/system/bcvk-journal-stream.service /run/systemd/system/; \
  cp /usr/lib/systemd/system/bcvk-journal-stream-initrd.service /run/systemd/system/; \
  ln -sf ../bcvk-journal-stream.service /run/systemd/system/sysinit.target.wants/; \
  ln -sf ../bcvk-journal-stream-initrd.service /run/systemd/system/sysinit.target.wants/; \
  systemctl --no-block start bcvk-journal-stream-initrd.service'
